# sys-apps/oomd/files/oomd.init
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="oomd"
description="userspace out-of-memory daemon"
command=/usr/bin/oomd
command_args="${oomd_args}"

depend() {
	need procfs sysfs
}
